# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file
 
[supervisord]
logfile=./supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=./supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=./supervisord/               ; where child log files will live
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:mesos_master]
command=/Users/palaitis/Developer/mesos/mesos-0.18.0/build/bin/mesos-master.sh --zk=zk://localhost:2181/mesos
stdout_logfile=./log/%(program_name)s.log
stderr_logfile=./log/%(program_name)s.log
autorestart=true

[program:mesos_slave]
command=/Users/palaitis/Developer/mesos/mesos-0.18.0/build/bin/mesos-slave.sh --master=zk://localhost:2181/mesos
stdout_logfile=./log/%(program_name)s.log
stderr_logfile=./log/%(program_name)s.log
autorestart=true

[program:marathon]
command=/Users/palaitis/Developer/mesos/marathon/bin/start --master zk://localhost:2181/mesos --zk zk://localhost:2181/marathon
stdout_logfile=./log/%(program_name)s.log
stderr_logfile=./log/%(program_name)s.log
autorestart=true

[program:zookeeper]
command=/Users/palaitis/Developer/mesos/zookeeper/zookeeper-3.4.6/bin/zkServer.sh start-foreground
stdout_logfile=./log/%(program_name)s.log
stderr_logfile=./log/%(program_name)s.log
autorestart=true
stopasgroup=true

